# data/session/libvirt-dbus.service.in
[Unit]
Description=Libvirt DBus Service

[Service]
BusName=org.libvirt
ExecStart=@sbindir@/libvirt-dbus --session

[Install]
Alias=org.libvirt.service
